# examples/launchd/com.carve.onewheel-ha-bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
  <dict>
    <key>Label</key>
    <string>com.carve.onewheel-ha-bridge</string>

    <key>ProgramArguments</key>
    <array>
      <string>/opt/onewheel-ha-bridge/.venv/bin/onewheel-ha-bridge</string>
      <string>--config</string>
      <string>/opt/onewheel-ha-bridge/config.toml</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/opt/onewheel-ha-bridge</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/opt/onewheel-ha-bridge/bridge.log</string>

    <key>StandardErrorPath</key>
    <string>/opt/onewheel-ha-bridge/bridge.log</string>
  </dict>
</plist>
